# server/base/backup/crontab
#minute hour    day     month   week  command
*/10       *       *       *       *     /cron-job.sh